# externals/supervisor/conf.d/sfroidsvr.conf
[program:sfroidsvr]
command = /home/sfroid/website/django/bin/gunicorn_start.sh
user = sfroiduser
stdout_logfile = /home/sfroid/website/django/logs/gunicorn_supervisor.log
redirect_stderr = true
